# odaf-pipeline.service
[Unit]
Description=ODAF Data Analytics Pipeline
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/odaf-pipeline
ExecStartPre=/bin/mkdir -p /opt/odaf-pipeline
ExecStartPre=/bin/chown -R root:root /opt/odaf-pipeline
ExecStart=/opt/odaf-pipeline/odaf_startup.sh
Restart=always
RestartSec=10
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=odaf-pipeline

[Install]
WantedBy=multi-user.target